# crontab
*/1 * * * * python3 /main.py >> /cron_log.log
